# docker/all-in-one/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/data/logs

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ====== MySQL ======
[program:mysql]
command=mysqld --user=mysql --datadir=/data/mysql --socket=/var/run/mysqld/mysqld.sock --port=3306 --character-set-server=utf8mb4 --collation-server=utf8mb4_general_ci --bind-address=127.0.0.1
autostart=true
autorestart=true
priority=10
stdout_logfile=/data/logs/mysql.log
stderr_logfile=/data/logs/mysql-error.log

# ====== Redis ======
[program:redis]
command=redis-server --port 6379 --requirepass %(ENV_REDIS_PASSWORD)s --dir /data/redis --appendonly yes --bind 127.0.0.1 --daemonize no
autostart=true
autorestart=true
priority=10
stdout_logfile=/data/logs/redis.log
stderr_logfile=/data/logs/redis-error.log

# ====== MinIO ======
[program:minio]
command=minio server /data/minio --address ":9000" --console-address ":9001"
environment=MINIO_ROOT_USER="%(ENV_MINIO_ACCESS_KEY)s",MINIO_ROOT_PASSWORD="%(ENV_MINIO_SECRET_KEY)s"
autostart=true
autorestart=true
priority=10
stdout_logfile=/data/logs/minio.log
stderr_logfile=/data/logs/minio-error.log

# ====== Java Backend ======
[program:server]
command=/app/start-server.sh
directory=/app
environment=MYSQL_PASSWORD="%(ENV_MYSQL_PASSWORD)s",MYSQL_DATABASE="%(ENV_MYSQL_DATABASE)s",REDIS_PASSWORD="%(ENV_REDIS_PASSWORD)s",MINIO_ACCESS_KEY="%(ENV_MINIO_ACCESS_KEY)s",MINIO_SECRET_KEY="%(ENV_MINIO_SECRET_KEY)s",MINIO_BUCKET="%(ENV_MINIO_BUCKET)s",JAVA_OPTS="%(ENV_JAVA_OPTS)s",S3_ENABLED="true",S3_PROVIDER="MINIO",MINIO_ENDPOINT="http://127.0.0.1:9000",MINIO_REGION="us-east-1"
autostart=true
autorestart=true
priority=20
startsecs=15
stdout_logfile=/data/logs/server.log
stderr_logfile=/data/logs/server-error.log

# ====== Nginx ======
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=30
stdout_logfile=/data/logs/nginx-access.log
stderr_logfile=/data/logs/nginx-error.log
